# big-first-boot/usr/lib/systemd/system/big-first-boot.service
[Unit]
Before=sddm.service

[Service]
ExecStart=/usr/bin/big-first-boot

[Install]
WantedBy=graphical.target
